# turbo/launchd/com.exo.claude.turbo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exo.claude.turbo</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/kitelev/Documents/exocortex-obsidian-plugin && exec bash turbo/watch.sh . "Auto-fix recent changes" 6 250 5 auto sonnet</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/kitelev/.npm-global/bin</string>
        <key>HOME</key>
        <string>/Users/kitelev</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/kitelev/Library/Logs/claude-turbo.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/kitelev/Library/Logs/claude-turbo.err</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/kitelev/Documents/exocortex-obsidian-plugin</string>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>